# lab-dashboard.service
[Unit]
Description=lab.local Homelab Dashboard
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=juyoungoh
Group=juyoungoh
SupplementaryGroups=docker
WorkingDirectory=/home/juyoungoh/nas/Projects/dashboard
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/juyoungoh/nas/Projects/dashboard/.venv/bin/python -m uvicorn backend.main:app --host 0.0.0.0 --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
